# webSlave.service
[Unit]
Description=TCP Server Service
After=network.target

[Service]
ExecStart=/usr/local/bin/webSlave
Restart=on-failure

[Install]
WantedBy=multi-user.target
